# backup/backup.service
[Unit]
Description=Nightly Backup Job

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/restic-backup.sh
